# examples/clear_cache.cron
33 22 * * * su - phoenix -c find /home/phoenix/anaconda/var/cache/pywps -name "*.nc" -atime +14 | xargs rm > /dev/null 2>&1
