# examples/clear_cache.cron
33 22 * * * su - phoenix -c 'find /home/phoenix/var/lib/cache/pywps -type f -name "*.nc" -atime +3 -print0 | xargs -0r rm'
33 23 * * * su - phoenix -c 'find /home/phoenix/var/lib/pywps/outputs -type f -ctime +14 -print0 | xargs -0r rm'

