# examples/clear_cache.cron
33 22 * * * su - phoenix -c find /home/phoenix/anacoda/var/cache/pywps -atime +14 | xargs rm -v > /dev/null 2>&1
